# deploy/systemd/airproche-cert-renew.timer
[Unit]
Description=Twice-daily Airproche certificate renewal check

[Timer]
OnCalendar=*-*-* 03,15:20:00 UTC
Persistent=true
RandomizedDelaySec=30m
Unit=airproche-cert-renew.service

[Install]
WantedBy=timers.target
